# Services/Database/db_email.service
[Unit]
Description="server for email related stuff"
After=mysql.service

[Service]
ExecStart=/home/phoenix/git/it490-sp26/sqlwork/emailServer.php
Restart=always
